# conf/supervisor_nerd2grip.ini
[program:grip_nerd2celery]
command=python3 /home/lenka/nerd2celery.py
directory=/home/lenka
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/grip/nerd2celery.err.log
stdout_logfile=/var/log/grip/nerd2celery.out.log
user=root
environment=SECRET_PASSPHRASE='this is secret',SECRET_TWO='another secret'
